# smb1-proxy/rootfs/etc/services.d/smbd/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "Stopping Samba daemon..."
pkill -TERM smbd 2>/dev/null || true
sleep 2
pkill -KILL smbd 2>/dev/null || true
