# smb1-proxy/rootfs/etc/services.d/smbd/run
#!/usr/bin/with-contenv sh
set -eu
mkdir -p /run/samba

ulimit -c 0 || true

is_alive() {
  _f="$1"
  if [ -f "$_f" ]; then
    _p="$(cat "$_f" || true)"
    [ -n "$_p" ] && kill -0 "$_p" 2>/dev/null && return 0
  fi
  return 1
}

if is_alive /run/samba/nmbd.pid; then
  echo "[info] nmbd already running (PID $(cat /run/samba/nmbd.pid)); not starting another."
else
  echo "[info] starting nmbd"
  rm -f /run/samba/nmbd.pid
  nmbd -F &
fi

echo "[info] starting smbd (foreground)"
exec smbd --foreground --no-process-group --debug-stdout --configfile=/etc/samba/smb-server.conf
