# afp2s3/rootfs/etc/services.d/minio/run
#!/usr/bin/with-contenv bash
set -euo pipefail

exec /usr/local/bin/minio server /mnt/afp --console-address ":9001"

