# smb1-proxy/rootfs/etc/services.d/dbus/run
#!/usr/bin/with-contenv sh
set -eu
mkdir -p /run/dbus
if [ -S /run/dbus/system_bus_socket ]; then
  if command -v dbus-send >/dev/null 2>&1 && dbus-send --system --dest=org.freedesktop.DBus / org.freedesktop.DBus.ListNames >/dev/null 2>&1; then
    echo "[info] System D-Bus already running; not starting a second instance."
    exec tail -f /dev/null
  fi
fi
rm -f /run/dbus/pid || true
exec /usr/bin/dbus-daemon --system --nofork --nopidfile
